# packaging/appliance/debian-wayland/bte-kiosk.service
[Unit]
Description=BTE Wayland Kiosk
Documentation=man:cage(1) man:foot(1)
After=systemd-user-sessions.service systemd-udevd.service network-online.target
Wants=network-online.target
Conflicts=getty@tty1.service

[Service]
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
StandardInput=tty
StandardOutput=journal+console
StandardError=journal+console
WorkingDirectory=/var/lib/bte
RuntimeDirectory=bte
RuntimeDirectoryMode=0700
Environment=XDG_RUNTIME_DIR=/run/bte
Environment=XDG_SESSION_TYPE=wayland
Environment=XDG_CURRENT_DESKTOP=bte
Environment=WLR_NO_HARDWARE_CURSORS=1
Environment=WLR_RENDERER=pixman
ExecStart=/usr/bin/cage -s -- /usr/bin/foot --config=/etc/bte/foot.ini /usr/local/bin/bte-launch
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
